# install/supervisor/io.solomon.hermes-gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.solomon.hermes-gateway</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>hermes</string>
        <string>gateway</string>
        <string>start</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/solomon/hermes-gateway.out.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/solomon/hermes-gateway.err.log</string>
</dict>
</plist>
