# install/supervisor/io.solomon.worker.template.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.solomon.worker.__SLUG__</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON__</string>
        <string>-m</string>
        <string>solomon.workers.__MODULE__</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__SOLOMON_ROOT__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>__SOLOMON_ROOT__</string>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/solomon/__SLUG__.out.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/solomon/__SLUG__.err.log</string>
</dict>
</plist>
